# templates/ansible/exasol-gre-mesh.service
[Unit]
Description=Ensure GRE mesh overlay for Exasol cluster
After=network-online.target systemd-networkd-wait-online.service
Wants=network-online.target systemd-networkd-wait-online.service
Before=c4_cloud_command.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/exasol-gre-mesh.sh
RemainAfterExit=yes
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
