# templates/ansible/exasol-data-symlinks.service
[Unit]
Description=Create stable Exasol data volume symlinks
Documentation=https://github.com/exasol/exasol-deployer
# Run after network is online (for metadata access) and local filesystems are mounted
After=network-online.target local-fs.target
# Ensure udev has finished processing devices
After=systemd-udev-settle.service
# For iSCSI-based providers (OCI): Wait for iSCSI volumes to be attached
After=open-iscsi.service
Wants=open-iscsi.service
# Run before any services that might need these symlinks
Before=docker.service
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/local/bin/exasol-data-symlinks.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

# Security hardening
ProtectSystem=strict
ReadWritePaths=/dev /etc/iscsi /var/lib/iscsi /run/lock
StateDirectory=iscsi
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
