# templates/ansible/exasol-overlay.service
[Unit]
Description=VXLAN overlay network for Exasol cluster multicast support
After=network-online.target
Wants=network-online.target
Before=exasol.service c4_cloud_command.service
Documentation=https://www.kernel.org/doc/Documentation/networking/vxlan.txt

[Service]
Type=oneshot
ExecStart=/usr/local/bin/exasol-overlay-setup.sh
ExecStop=/bin/bash -c 'ip link show vxlan-br0 >/dev/null 2>&1 && ip link delete vxlan-br0 || true; ip link show vxlan0 >/dev/null 2>&1 && ip link delete vxlan0 || true'
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
