# templates/ansible/exasol-overlay.service
[Unit]
Description=VXLAN overlay network for Exasol cluster multicast support
After=network-online.target
Wants=network-online.target
Before=exasol.service c4_cloud_command.service
Documentation=https://www.kernel.org/doc/Documentation/networking/vxlan.txt

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/exasol-overlay-setup.sh
ExecStop=/usr/local/bin/exasol-overlay-stop.sh
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
TimeoutStartSec=60s
TimeoutStopSec=30s

[Install]
WantedBy=multi-user.target
